# data/bin/JobAdderServer.service
[Unit]
Description=JobAdder Server Daemon
Requires=postgresql.service
After=postgresql.service
AssertPathExists=

[Service]
ExecStart=/usr/bin/ja-server
ExecStop=/usr/bin/bash -c "/usr/bin/ja-server -k && tail --pid=$MAINPID -f /dev/null"

[Install]
WantedBy=multi-user.target
