# data/bin/JobAdderServer.service
[Unit]
Description=JobAdder Server Daemon
Requires=postgresql.service
After=postgresql.service
AssertPathExists=

[Service]
ExecStart=/usr/bin/ja-server
# Wait 2s for the server to stop.
ExecStop=/usr/bin/bash -c '/usr/bin/ja-server -k && sleep 2'

[Install]
WantedBy=multi-user.target
